# scripts/systemd/lsats-bronze-ad.service
[Unit]
Description=LSATS Bronze Ingestion - Active Directory
After=postgresql.service

[Service]
Type=oneshot
User=lsats
WorkingDirectory=/opt/LSATS_Data_Hub
EnvironmentFile=/etc/LSATS_Data_Hub/hub.conf
Environment=PGPASSFILE=/opt/LSATS_Data_Hub/.pgpass
# Decrypt AD and MCommunity LDAP passwords at service start.
# Credentials are stored as machine-encrypted blobs in /etc/credstore/ (systemd-creds).
# The orchestrator script reads $CREDENTIALS_DIRECTORY/{ad_password,ldap_password}
# and exports them as AD_PASSWORD and LDAP_PASSWORD before invoking Python.
LoadCredential=ad_password:/etc/credstore/ad_password.cred
LoadCredential=ldap_password:/etc/credstore/ldap_password.cred
ExecStart=/bin/bash /opt/LSATS_Data_Hub/scripts/database/orchestrate_bronze_ad.sh
StandardOutput=journal
StandardError=journal
